# config/systemd/sauna-ctrl-deploy.timer
[Unit]
Description=Poll GitHub for sauna_ctrl updates every 30s

[Timer]
OnBootSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
